# services/RadioTemperature/run
#!/bin/sh
exec 2>&1
exec /data/RadioTemperature/RadioTemperature.py
